# infrastructure/src/crontab/umigratecron
0 4 * * * export $(egrep -v '^#' /home/umigrate/venv/.env | xargs) && python3 /home/umigrate/api/manage.py clearsessions
