# infrastructure/src/crontab/umigratecron
0 4 * * * python3 /home/umigrate/api/manage.py clearsessions
